# deploy/twin-roll.service
[Unit]
Description=Twin Roll Ludo Game Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/twin-roll
ExecStart=/opt/twin-roll/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
